# workers/aws/prefect-aws-worker-watch.timer
# /etc/systemd/system/prefect-aws-worker-watch.timer
#
# Fires every 2 minutes to run prefect-aws-worker-watch.service, which re-runs
# `docker compose up -d` to restore any stopped/removed worker containers.

[Unit]
Description=Run prefect-aws-worker-watch every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=30s
Unit=prefect-aws-worker-watch.service

[Install]
WantedBy=timers.target
